# deploy/supervisor.conf
[program:tms]
command=/var/www/tms/venv/bin/gunicorn -c deploy/gunicorn.conf wsgi:app
directory=/var/www/tms
user=root
autostart=true
autorestart=true
environment = MODE="PRODUCTION"